# board/kiosk/rootfs_overlay/etc/inittab
# /etc/inittab for Kiosk OS (BusyBox init)
#
# Format: <id>:<runlevels>:<action>:<process>
# Note: BusyBox init ignores runlevels, only id:action:process matter

# System startup - mount filesystems and run init scripts
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t devtmpfs devtmpfs /dev
::sysinit:/bin/mkdir -p /dev/pts /dev/shm
::sysinit:/bin/mount -t devpts devpts /dev/pts
::sysinit:/bin/mount -t tmpfs tmpfs /dev/shm
::sysinit:/bin/mount -t tmpfs tmpfs /tmp
::sysinit:/bin/mount -t tmpfs tmpfs /run
::sysinit:/bin/ln -sf /run /var/run
::sysinit:/bin/ln -sf /tmp /var/tmp
::sysinit:/bin/mkdir -p /var/lock /var/log
::sysinit:/bin/hostname -F /etc/hostname

# Run all init scripts (networking, dbus, etc.)
::sysinit:/etc/init.d/rcS

# Start kiosk on tty1 after init completes
# Using 'once' instead of 'respawn' to avoid restart loops
::once:/usr/bin/kiosk-start

# Emergency shell on tty2 (Ctrl+Alt+F2)
tty2::askfirst:/bin/sh

# Shutdown hooks
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::shutdown:/bin/umount -a -r
